# public-release/scripts/security/security-monitor.service
[Unit]
Description=Security Event Monitor and Automated Response
After=network.target docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/scripts/security
ExecStart=/usr/bin/python3 /opt/scripts/security/event-monitor.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=false  # Needs root for iptables
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=false  # Needs to write logs

# Resource limits
CPUQuota=50%
MemoryLimit=512M

[Install]
WantedBy=multi-user.target