# scripts/security/ssh-monitor.service
[Unit]
Description=SSH Login Monitor
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/env bash -c 'source /etc/profile && /opt/scripts/security/ssh-monitor.sh'
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
